# src/meta-product/recipes-application/alive-service/files/systemd/alive.service
[Unit]
Description=Alive Service
After=network.target

[Service]
Type=simple
ExecStart=/srv/alive-service
Restart=no

[Install]
WantedBy=multi-user.target
